# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/config/logs/supervisor

# ==================================================================
# 1. xrdp 服务 (Root)
# ==================================================================
[program:xrdp]
command=/usr/sbin/xrdp -nodaemon
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/config/logs/xrdp.log

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/config/logs/xrdp-sesman.log

# ==================================================================
# 2. Xvfb (虚拟显示器) - 必须先启动
# 注意：这里以 app 用户启动，以便 x11vnc 可以访问 .Xauthority
# ==================================================================
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac
autostart=true
autorestart=true
user=%(ENV_USER_NAME)s
environment=DISPLAY=":99", HOME="/config"
priority=10
stdout_logfile=/config/logs/xvfb.log

# ==================================================================
# 3. 桌面环境 (XFCE + D-Bus)
# 依赖 Xvfb，必须由同一个用户 (app) 启动
# ==================================================================
[program:xfce]
# 关键修正：
# 1. 等待 Xvfb 启动 (sleep 2)
# 2. 使用 dbus-launch 包裹
# 3. 移除旧的 xfdesktop 和 冲突的 xfce program
command=/bin/sh -c "sleep 2 && export DISPLAY=:99 && /usr/bin/dbus-launch --exit-with-session /usr/bin/xfce4-session"
autostart=true
autorestart=true
user=%(ENV_USER_NAME)s
environment=DISPLAY=":99", HOME="/config", DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
priority=20
stdout_logfile=/config/logs/xfce.log
stderr_logfile=/config/logs/xfce.err

# ==================================================================
# 4. x11vnc (AppUser)
# ==================================================================
[program:x11vnc]
command=/usr/bin/x11vnc -forever -shared -rfbport 5901 -display :99 -rfbauth /tmp/vnc.passwd -o /config/logs/x11vnc.log
autostart=true
autorestart=true
user=%(ENV_USER_NAME)s
environment=DISPLAY=":99", HOME="/config"
priority=30
redirect_stderr=true

# ==================================================================
# 5. noVNC Web 服务
# ==================================================================
[program:novnc]
command=/usr/share/noVNC/utils/novnc_proxy --vnc localhost:5901 --web /usr/share/noVNC
autostart=true
autorestart=true
user=%(ENV_USER_NAME)s
environment=HOME="/config"
priority=40
redirect_stderr=true
stdout_logfile=/config/logs/novnc.log